# agent-articole.service
[Install]
WantedBy=multi-user.target

[Unit]                                                                                                                                          
Description=Agent Articole FastAPI
After=network.target

[Service]
Type=exec
User=ubuntu
WorkingDirectory=/home/ubuntu/claude/agent_articole
ExecStart=/home/ubuntu/claude/agent_articole/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8002 --log-level info
Restart=on-failure
RestartSec=5                                                                                                                                    
                  
[Install]
WantedBy=multi-user.target
